# repository/autostart-x11/.slash/etc/profile.d/999-autostart-X11-on-first-tty.sh
if [ -z "$DISPLAY" ] && [ -n "$XDG_VTNR" ]; then
  exec startx
fi
