# repository/tmuxifier/.slash/etc/profile.d/100-tmuxifier.sh
export PATH="$PATH:$HOME/.slash/lib/tmuxifier/bin"
export TMUXIFIER_LAYOUT_PATH="$HOME/.slash/etc/tmuxifier/"
eval "$(tmuxifier init -)"
