# repository/tmux/.slash/etc/profile.d/900-tmux.sh
# load bashrc on tmux session
case $- in *i*) . ~/.bashrc;; esac
