# common/nexway/.slash/etc/profile.d/nexway-ssh.sh
_nexway-ssh()
{
    local cur prev opts
    COMPREPLY=()
    cur="${COMP_WORDS[COMP_CWORD]}"
    prev="${COMP_WORDS[COMP_CWORD-1]}"
    opts="softgallery@devmrogacki softgallery N3@tenwadev N3@tenwaprep --help -h"

    COMPREPLY=( $(compgen -W "${opts}" -- ${cur}) )
    return 0
}
complete -F _nexway-ssh nexway-ssh
