# common/urxvt/.slash/etc/profile.d/100-urxvt.sh
export URXVT_PERL_LIB="$HOME"/.urxvt/ext/
xrdb -merge "$HOME"/.urxvt/conf
