# common/games/.slash/etc/profile.d/tibia.sh
export PATH=$PATH:$HOME/.slash/usr/bin/tibiaprices
